# systemd/backup-sh.timer
[Timer]
OnCalendar=*:00/15
AccuracySec=1min
RandomizedDelaySec=1min
Persistent=yes

[Install]
WantedBy=timers.target
